# systemd/sid.service.in
#
# SPDX-FileCopyrightText: (C) 2017-2024 Red Hat, Inc.
#
# SPDX-License-Identifier: GPL-2.0-or-later
#

[Unit]
Description=Storage Instantiation Daemon
DefaultDependencies=no
After=systemd-udevd.service
Before=sysinit.target

[Service]
Type=notify
EnvironmentFile=(SYSCONFIGDIR)/sid.sysconfig
Environment=SERVICE_ACTIVATION_TYPE=FD_PRELOAD
ExecStart=(SBINDIR)/sid -f
