# systemd/sid.service.in
[Unit]
Description=Storage Instantiation Daemon
DefaultDependencies=no
After=systemd-udevd.service
Before=sysinit.target

[Service]
Type=notify
EnvironmentFile=(SYSCONFIGDIR)/sid.sysconfig
Environment=SERVICE_ACTIVATION_TYPE=FD_PRELOAD
ExecStart=(BINDIR)/sid -f
