# systemd/sid.socket.in
#
# SPDX-FileCopyrightText: (C) 2017-2025 Red Hat, Inc.
#
# SPDX-License-Identifier: GPL-2.0-or-later
#

[Unit]
Description=Storage Instantiation Daemon socket
DefaultDependencies=no
Before=sockets.target
After=systemd-udevd-kernel.socket

[Socket]
ListenStream=@sid-ubridge.socket
PassCredentials=yes

[Install]
WantedBy=sockets.target
