# systemd/sid.socket.in
[Unit]
Description=Storage Instantiation Daemon socket
DefaultDependencies=no
Before=sockets.target
After=systemd-udevd-kernel.socket

[Socket]
ListenStream=@sid-ubridge.socket
PassCredentials=yes

[Install]
WantedBy=sockets.target
